# packages/linux-firmware/brcmfmac_sdio-firmware-aml/system.d/brcmfmac_sdio-firmware@.service
[Unit]
Description=Broadcom sdio firmware update for %I
ConditionPathExists=/dev/ttyS1

[Service]
Type=simple
RemainAfterExit=yes
ExecStart=/bin/sh -c "/sbin/rfkill bluetooth; /usr/bin/brcm_patchram_plus --patchram /lib/firmware/brcm/%I.hcd --no2bytes --baudrate 2000000 --use_baudrate_for_download /dev/ttyS1; /usr/bin/hciattach /dev/ttyS1 any 2000000"
