# packages/mediacenter/ExtroCore/system.d/custom_backup-restore.target
[Unit]
Description=Backup and restore
Requires=backup-restore.service
After=backup-restore.service
AllowIsolate=yes
