# packages/mediacenter/ExtroCore/system.d/custom_backup-restore.service
[Unit]
Description=Backup and restore
DefaultDependencies=no

[Service]
Type=idle
ExecStart=/usr/lib/libreelec/custom_backup-restore
StandardInput=tty-force
StandardOutput=inherit
StandardError=inherit
