# packages/linux-drivers/bcmdhd-sdio-aml/system.d/bcmdhd-sdio-aml-suspend.service
[Unit]
Description=Broadcom WiFi suspend service
Before=sleep.target

[Service]
Type=oneshot
ExecStart=/sbin/ifconfig wlan0 down
TimeoutSec=0

[Install]
WantedBy=sleep.target
