# packages/linux-drivers/bcmdhd-sdio-aml/system.d/bcmdhd-sdio-aml-resume.service
[Unit]
Description=Broadcom WiFi resume service
After=suspend.target

[Service]
Type=oneshot
ExecStart=/sbin/ifconfig wlan0 up
TimeoutSec=0

[Install]
WantedBy=suspend.target
